# deploy/systemd/soc-sentinel.service
[Unit]
Description=SOC Sentinel production service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=soc-sentinel
Group=soc-sentinel
WorkingDirectory=/opt/soc-sentinel/SOC-Sentinel/soc_server
EnvironmentFile=/opt/soc-sentinel/SOC-Sentinel/.env
ExecStart=/opt/soc-sentinel/SOC-Sentinel/.venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 --access-logfile - --error-logfile - --timeout 120 wsgi:application
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
