# deploy/systemd/soc-sentinel.service
[Unit]
Description=SOC Sentinel production service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=soc-sentinel
Group=soc-sentinel
WorkingDirectory=__SOC_SENTINEL_HOME__/soc_server
EnvironmentFile=__SOC_SENTINEL_HOME__/.env
ExecStart=__SOC_SENTINEL_HOME__/.venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 --access-logfile - --error-logfile - --timeout 120 wsgi:application
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
